# supervisor_snapraid.ini
[supervisord]
nodaemon=true

[program:cron]
startsecs=0
autostart=true
autorestart=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
command=start_cron.sh

[program:snapraid_runner]
startsecs=0
autostart=false
autorestart=false
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
command=snapraid-runner.py -c /config/snapraid-runner.conf
